# ipac-ng.cron
*/5 * * * *	root umask 022; /usr/sbin/fetchipac

